# control-plane/agents/scribe/scribe.service
[Unit]
Description=THE SCRIBE - Council Meeting Secretary
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/leveredge/control-plane/agents/scribe
Environment=PYTHONPATH=/opt/leveredge/control-plane
ExecStart=/usr/bin/python3 -m uvicorn scribe:app --host 0.0.0.0 --port 8301
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
