# control-plane/agents/hades/hades.service
[Unit]
Description=HADES Rollback System Agent
After=network.target event-bus.service chronos.service

[Service]
Type=simple
User=damon
WorkingDirectory=/opt/leveredge/control-plane/agents/hades
Environment="N8N_CONTROL_URL=https://control.n8n.leveredgeai.com"
Environment="N8N_PROD_URL=https://n8n.leveredgeai.com"
Environment="N8N_USER=admin"
Environment="N8N_PASS=oMtnAe9qsrxhMe/1ROmokg=="
Environment="EVENT_BUS_URL=http://localhost:8099"
Environment="CHRONOS_URL=http://localhost:8010"
ExecStart=/home/damon/.local/bin/uvicorn hades:app --host 0.0.0.0 --port 8008
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
