# shared/systemd/leveredge-chronos.service
[Unit]
Description=LeverEdge CHRONOS - Time Lord of Backups
After=network.target
Wants=leveredge-event-bus.service leveredge-hermes.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/leveredge/control-plane/agents/chronos
Environment=PYTHONUNBUFFERED=1
Environment=REGISTRY_PATH=/opt/leveredge/config/agent-registry.yaml
Environment=EVENT_BUS_URL=http://localhost:8099
Environment=HERMES_URL=http://localhost:8014
ExecStart=/opt/leveredge/shared/venv/bin/python chronos.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
