# shared/systemd/leveredge-hermes.service
[Unit]
Description=LeverEdge HERMES - Divine Messenger
After=network.target
Wants=leveredge-event-bus.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/leveredge/control-plane/agents/hermes
Environment=PYTHONUNBUFFERED=1
Environment=REGISTRY_PATH=/opt/leveredge/config/agent-registry.yaml
Environment=EVENT_BUS_URL=http://localhost:8099
ExecStart=/opt/leveredge/shared/venv/bin/python hermes.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
