# shared/systemd/schema-sync.timer
[Unit]
Description=Nightly schema sync DEV to PROD

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
