# shared/systemd/agent-services/clio.service
[Unit]
Description=CLIO - History & Documentation Agent
Documentation=https://github.com/leveredge/creative-fleet
After=network.target event-bus.service
Wants=event-bus.service

[Service]
Type=simple
User=damon
Group=damon
WorkingDirectory=/opt/leveredge/control-plane/agents/clio
Environment="PATH=/usr/local/bin:/usr/bin:/bin:/home/damon/.local/bin"
Environment="HOME=/home/damon"
Environment="EVENT_BUS_URL=http://localhost:8099"
ExecStart=/home/damon/.local/bin/uvicorn clio:app --host 0.0.0.0 --port 8021
Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=clio

[Install]
WantedBy=multi-user.target
