# shared/systemd/agent-services/meal-planner.service
[Unit]
Description=MEAL-PLANNER - Meal Planning & Recipe Agent
Documentation=https://github.com/leveredge/meal-planner
After=network.target event-bus.service
Wants=event-bus.service

[Service]
Type=simple
User=damon
Group=damon
WorkingDirectory=/opt/leveredge/control-plane/agents/meal-planner
Environment="PATH=/usr/local/bin:/usr/bin:/bin:/home/damon/.local/bin"
Environment="HOME=/home/damon"
Environment="EVENT_BUS_URL=http://localhost:8099"
ExecStart=/home/damon/.local/bin/uvicorn meal_planner:app --host 0.0.0.0 --port 8029
Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=meal-planner

[Install]
WantedBy=multi-user.target
