# shared/systemd/agent-services/nutritionist.service
[Unit]
Description=NUTRITIONIST - Nutrition & Diet Agent
Documentation=https://github.com/leveredge/nutritionist
After=network.target event-bus.service
Wants=event-bus.service

[Service]
Type=simple
User=damon
Group=damon
WorkingDirectory=/opt/leveredge/control-plane/agents/nutritionist
Environment="PATH=/usr/local/bin:/usr/bin:/bin:/home/damon/.local/bin"
Environment="HOME=/home/damon"
Environment="EVENT_BUS_URL=http://localhost:8099"
ExecStart=/home/damon/.local/bin/uvicorn nutritionist:app --host 0.0.0.0 --port 8028
Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=nutritionist

[Install]
WantedBy=multi-user.target
