# shared/systemd/agent-services/varys.service
[Unit]
Description=VARYS - Intelligence & Information Gathering Agent
Documentation=https://github.com/leveredge/varys
After=network.target event-bus.service
Wants=event-bus.service

[Service]
Type=simple
User=damon
Group=damon
WorkingDirectory=/opt/leveredge/control-plane/agents/varys
Environment="PATH=/usr/local/bin:/usr/bin:/bin:/home/damon/.local/bin"
Environment="HOME=/home/damon"
Environment="EVENT_BUS_URL=http://localhost:8099"
ExecStart=/home/damon/.local/bin/uvicorn varys:app --host 0.0.0.0 --port 8040
Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=varys

[Install]
WantedBy=multi-user.target
